# php-fpm/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * php /code/artisan schedule:run >> /proc/1/fd/1 2>/proc/1/fd/2
